# systemd/embedding.service
[Unit]
Description=Embedding Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/Mentara/microservices/quick-embedding
ExecStart=/root/Mentara/microservices/quick-embedding/build/bin/embedding_service
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target